# etc/network/if-up.d/sinc-puppet
/usr/share/linex-ubuntu-puppet/sinc-puppet